# docker/s6-overlay/s6-rc.d/addon-migration/up
#!/command/execlineb -P

# Run addon migrations after Coolify's main migrations
with-contenv
cd /var/www/html
foreground {
    php
    artisan
    migrate
    --path=vendor/amirhmoradi/coolify-enhanced/database/migrations
    --force
}

# Ensure custom templates cache directory is writable by www-data
# Host path /data/coolify/custom-templates is mounted to this container path
foreground { mkdir -p /var/www/html/storage/app/custom-templates }
foreground { chown www-data:www-data /var/www/html/storage/app/custom-templates }
